# chef/site-cookbooks/gitorious/templates/default/gitorious.init.erb
#! /bin/sh
### BEGIN INIT INFO
# Provides:          gitorious
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start a Rails instance
# Description:       Do the simplest thing possible in keeping with
#             upstart to spin up a single Rails instance.
### END INIT INFO

# Author: Sam Pointer 

PATH=/sbin:/usr/sbin:/bin:/usr/bin
USER="git"
PORT=3000
RAILS_ROOT="/var/www/git.myserver.com/gitorious"
COMMAND="rvm ree ruby script/server -e production -p $PORT -d"
DESCRIPTION="Gitorious"

# Load the VERBOSE setting and other rcS variables
. /lib/init/vars.sh

# Define LSB log_* functions.
# Depend on lsb-base (>= 3.0-6) to ensure that this file is present.
. /lib/lsb/init-functions

#
# Function that starts the daemon/service
#
do_start()
{
    # Return
    #   0 if daemon has been started
    #   1 if daemon was already running
    #   2 if daemon could not be started
    su - $USER -c "cd $RAILS_ROOT && $COMMAND"
}

case "$1" in
  start)
        [ "$VERBOSE" != no ] && log_daemon_msg "Starting $DESCRIPTION"
        do_start
        case "$?" in
                0|1) [ "$VERBOSE" != no ] && log_end_msg 0 ;;
                2) [ "$VERBOSE" != no ] && log_end_msg 1 ;;
        esac
        ;;
esac
